# etc/script/systemd/bmstools.service
[Unit]
Description=bmstools service

[Service]
Restart=always
User=root
ExecStart=/usr/bin/bmstools

[Install]
WantedBy=multi-user.target